# profile.d/d3d12.sh
#!/bin/bash

export LIBVA_DRIVER_NAME=d3d12
export MESA_D3D12_DEFAULT_ADAPTER_NAME=nvidia