# profile.d/http_proxy_env.sh
export HTTPS_PROXY=http://gateway.zscloud.net:10268
export HTTP_PROXY=http://gateway.zscloud.net:10268
export https_proxy=http://gateway.zscloud.net:10268
export http_proxy=http://gateway.zscloud.net:10268
export NO_PROXY="localhost,.mediahub.ch,.swisstxt.ch,dev-camunda.mediahub.ch,stage-camunda.mediahub.ch,camunda.mediahub.ch,$NO_PROXY"
