# profile.d/vault.sh
#!/bin/bash

export VAULT_ADDR="https://vault.w084-prd.aws-ec1.srgssr.cloud"
export VAULT_FORMAT=json
