# profile.d/vault.sh
#!/bin/bash

export VAULT_ADDR=https://vault.swisstxt.ch
export VAULT_FORMAT=json