# profile.d/wslu_browser.sh
#!/bin/bash

# launch wsl host browser
export BROWSER=wslview